# services/weather-dashboard.service
[Unit]
Description=Weather Station Dashboard
After=network.target influxdb.service

[Service]
Type=simple
User=masterbox
Group=masterbox
WorkingDirectory=/home/masterbox/weatherstation
ExecStart=/home/masterbox/weatherstation/.venv/bin/gunicorn --config config/gunicorn_config.py wsgi:application
Restart=always
RestartSec=10
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target